# server-config/quadlets/pocket-surf-vikunja.container
[Unit]
Description=Pocket Surf Vikunja (Task Management)
After=pocket-surf.pod.service
Requires=pocket-surf.pod.service

[Container]
Image=docker.io/vikunja/vikunja:latest
Pod=pocket-surf.pod

# Persistent storage for DB and uploaded files
Volume=vikunja-files:/app/vikunja/files:Z

# Use SQLite – no separate database container needed
Environment=VIKUNJA_SERVICE_PUBLICURL=http://tasks.local
Environment=VIKUNJA_DATABASE_TYPE=sqlite
Environment=VIKUNJA_DATABASE_PATH=/app/vikunja/files/vikunja.db

Restart=always
RestartSec=10s

HealthCmd=/usr/bin/wget --no-verbose --tries=1 --spider http://localhost:3456/api/v1/info || exit 1
HealthInterval=30s
HealthRetries=3
HealthStartPeriod=30s

[Service]
Restart=always
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target default.target
